# etc/service/php-fpm/run
#!/bin/bash
exec php-fpm --nodaemonize -R
